# packaging/prerm
#!/bin/sh
set -e
if [ "$1" = "remove" ]; then
    systemctl stop t2
    systemctl disable t2
fi
